# systemd/auth-monitor.service
[Unit]
Description=AuthMonitor

[Service]
ExecStart=/usr/bin/auth-monitor-daemon.sh --test
Restart=on-failure

[Install]
WantedBy=multi-user.target
